# scripts/com.todoist-things.bidirectional.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.todoist-things.bidirectional</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__SCRIPT_PATH__</string>
    </array>
    
    <key>StartInterval</key>
    <integer>300</integer><!-- Sync every 5 minutes -->
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/todoist-things-sync.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/todoist-things-sync.stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>TODOIST_THINGS_WORKER_URL</key>
        <string>__WORKER_URL__</string>
    </dict>
</dict>
</plist>